# systemd/thricegrip-stream.service
[Unit]
Description=ThriceGrip Video Stream (ustreamer)
After=thricegrip-gadget.service
Wants=thricegrip-gadget.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ustreamer \
    --device=/dev/video0 \
    --resolution=1920x1080 \
    --format=uyvy \
    --encoder=hw \
    --quality=80 \
    --host=0.0.0.0 \
    --port=8080
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
